# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

[program:bot]
directory=/app/bot
command=python main.py
autostart=true
autorestart=true

[program:miniapp]
directory=/app/mini-app/backend
command=node index.js
autostart=true
autorestart=true